# supervisor.conf
[program:laravel_queue]
command=/usr/local/bin/php /var/www/html/artisan queue:work
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/proc/1/fd/1
stderr_logfile=/proc/1/fd/2
